# wb-stocks/deploy/gcp/systemd/wb-forecast-recalc.service
[Unit]
Description=wb-stocks: recalculate sales forecast (MVP)
After=network-online.target wb-import-stocks.service wb-update-supplies.service
Wants=network-online.target

[Service]
Type=oneshot
User=wbstocks
Group=wbstocks
WorkingDirectory=/srv/wb-stocks/current
TimeoutStartSec=60min
ExecStart=/usr/bin/node --env-file=/etc/wb-stocks.env --import tsx scripts/run-sales-forecast-mvp.ts
Nice=15

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/wb-stocks/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

StandardOutput=journal
StandardError=journal
