# wb-stocks/deploy/gcp/systemd/wb-update-supplies.timer
[Unit]
Description=wb-stocks: daily WB supplies sync

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=10min
Unit=wb-update-supplies.service

[Install]
WantedBy=timers.target
